# ambientLightFix.service
[Unit]
Description=Automatically adjust screen-brightness

[Service]
Type=simple
ExecStart=/usr/bin/ambientLightFix

[Install]
WantedBy=multi-user.target
